# deploy/systemd/seapedia-backend.service
[Unit]
Description=SEAPEDIA Backend API
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/seapedia/backend
EnvironmentFile=/opt/seapedia/backend/.env
ExecStart=/opt/seapedia/backend/bin/seapedia
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
