# deploy/systemd/seapedia-frontend.service
[Unit]
Description=SEAPEDIA Frontend (Next.js)
After=network.target seapedia-backend.service

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/seapedia/frontend
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/home/ec2-user/.bun/bin/bun run start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
